# runner-clean.service
[Unit]
Description=Gitlab Runner Cache Clean

[Service]
Type=oneshot

ExecStart=/usr/share/gitlab-runner/clear-docker-cache

