# runner-clean.timer
[Unit]
Description=Runner Clean Timer

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
