# docker.service.d/override.conf
[Unit]

# require the mount /var/lob/docker/volumes

Requires=var-lib-docker-volumes.mount
BindsTo=var-lib-docker-volumes.mount
After=var-lib-docker-volumes.mount